# .config/autostart/xremap.desktop
[Desktop Entry]
Name=Xremap
Type=Application
Exec=xremap /home/takahisa/.config/xremap/config.yml --watch
Terminal=false
